# etc/init.d/rpi-surveillance
#!/bin/sh
### BEGIN INIT INFO
# Provides:          rpi-surveillance
# Required-Start:    $syslog
# Required-Stop:     $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start RPi Surveillance server at boot time
# Description:       Enable service provided by daemon.
### END INIT INFO

rpi_surveillance_path=/home/pi/RaspiSecurity
rpi_surveillance_pidfile=/var/run/rpi-surveillance.pid
rpi_camera_restart=/root/rpi-camera-restart
rpi_surveillance_name=$(basename "$0")

case "$1" in
    start)
        echo "Starting ${rpi_surveillance_name}"
        if [ -f "${rpi_surveillance_pidfile}" ]; then
            "./$0" restart
        else
            cd "${rpi_surveillance_path}" && gunicorn3 server:app -b 0.0.0.0:5555 -p "${rpi_surveillance_pidfile}" -D
        fi
        if [ -f "${rpi_camera_restart}" ]; then
            sleep 10
            curl -4 -s http://localhost:5555/start
            rm -f "${rpi_camera_restart}"
        fi
        ;;
    stop)
        echo "Stopping ${rpi_surveillance_name}"
        if [ -f "${rpi_surveillance_pidfile}" ]; then
            kill "$(cat "${rpi_surveillance_pidfile}")"
            rm "${rpi_surveillance_pidfile}"
        fi
        ;;
    restart)
        # Hacky way to restart
        status=$(curl -4 -s http://localhost:5555/status)
        "./$0" stop
        "./$0" start
        if [ "${status}" = Running ]; then
            sleep 2
            curl -4 -s http://localhost:5555/start > /dev/null
            exit $?
        fi
        ;;
    start-camera)
        echo "Starting ${rpi_surveillance_name} camera"
        curl -4 -s http://localhost:5555/start > /dev/null
        ;;
    stop-camera)
        echo "Stopping ${rpi_surveillance_name} camera"
        curl -4 -s http://localhost:5555/stop > /dev/null
        ;;
    status)
        if [ -f "${rpi_surveillance_pidfile}" ]; then
            status=$(curl -4 -s http://localhost:5555/status)
            echo "Server running; camera ${status}"
        else
            echo "Server stopped"
            exit 1
        fi
        ;;
    *)
        echo "Usage: $0 {start|stop|restart|start-camera|stop-camera|status}"
        exit 1
        ;;
esac

exit 0
